# data/export/upstart_user/process.conf.erb
start on starting <%= app %>-<%= name %>
stop on stopping <%= app %>-<%= name %>
respawn

chdir <%= engine.root %>
env PORT=<%= port %>
<% engine.env.each_pair do |var,env| %>
env <%= var.upcase %>=<%= shell_quote(env) %>
<% end %>

exec <%= process.command %> >> <%= log %>/<%=name%>-<%=num%>.log 2>&1
